# launchd/com.yourname.session-loader.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yourname.session-loader</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>backend.main:app</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8080</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/mac-session-loader</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>SESSION_LOADER_KEY</key>
        <string>your-secret-key-here</string>
        <key>DATA_DIR</key>
        <string>/Users/YOUR_USERNAME/mac-session-loader</string>
        <key>DEFAULT_WORKSPACE</key>
        <string>/Users/YOUR_USERNAME/</string>
        <key>DEFAULT_TIMEZONE</key>
        <string>America/New_York</string>
        <key>SESSION_IDLE_TIMEOUT_HOURS</key>
        <string>6.0</string>
        <key>SESSION_CLEANUP_INTERVAL_MINUTES</key>
        <string>10</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/session-loader.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/session-loader.log</string>
</dict>
</plist>
